# 2.4.125-ubuntu/files/etc/rc.local
#!/bin/sh -e
echo never > /sys/kernel/mm/transparent_hugepage/enabled

echo 1024 > /proc/sys/net/core/somaxconn

sysctl vm.overcommit_memory=1

exit 0